# service/umlaut.service
[Unit]
Description=Umlaut Daemon - System-wide keyboard remapping service
After=default.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/umlaut-scripts/umlaut_daemon.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
StartLimitBurst=5
StartLimitIntervalSec=30
TimeoutStopSec=5
KillMode=process
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
